# systemd/ndncert-ca.service.in
[Unit]
Description=Certificate Management Identity Management Service for NDN

[Service]
Environment=HOME=%S/ndncert-ca
ExecStart=@BINDIR@/ndncert-ca-server
Restart=on-failure
RestartPreventExitStatus=2
RestartSec=5
User=ndn

[Install]
WantedBy=multi-user.target
Alias=ndncert.service