# systemd/ndncert-ca.service.in
[Unit]
Description=Certificate Management Identity Management Service for NDN

[Service]
Environment=HOME=%S/ndncert
EnvironmentFile=-@SYSCONFDIR@/ndncert
EnvironmentFile=-@SYSCONFDIR@/default/ndncert
ExecStart=@BINDIR@/ndncert-ca-server
Restart=on-failure
RestartPreventExitStatus=2
User=ndn

[Install]
WantedBy=multi-user.target
